# alpine-node-red/root/etc/services.d/avahi-daemon/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}

umask "$UMASK_SET"

/usr/sbin/avahi-daemon -f /config/avahi.conf
